# tvheadend/rootfs/etc/services.d/tvheadend/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts TVHeadend
# ==============================================================================

bashio::log.info "Starting TVHeadend..."

exec /usr/bin/tvheadend -u root -g root -c /config/tvheadend --satip_xml http://192.168.1.1:49000/satipdesc.xml
